# files/etc/crontabs/root
0 4 * * * /etc/coremark.sh
0 6 * * 1 /usr/share/openclash/openclash_rule.sh
0 7 * * 1 /usr/share/openclash/openclash_ipdb.sh
0 7 * * 1 /usr/share/openclash/openclash_geosite.sh
0 7 * * 1 /usr/share/openclash/openclash_geoip.sh
0 7 * * 1 /usr/share/openclash/openclash_geoasn.sh
0 7 * * 1 /usr/share/openclash/openclash_chnroute.sh
0 12 * * * /etc/init.d/openclash restart
*/6 * * * * /etc/init.d/frpgl start 2>/dev/null
